# simple-docker-app.service
[Unit]
Description=Simple Docker App - Node.js Web Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/simple-docker-app
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=10

# Environment variables
Environment=NODE_ENV=production
Environment=PORT=3000

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=simple-docker-app

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
